# Fichiers/~ Polytech/NVIDIA/TD-9-11-22/Tanguy/dump-database.timer
[Unit]
Description=lancement toutes les 10 min

[Timer]
OnCalendar=*:0/10
AccuracySec=1s
Unit=backup.Service

[Install]
WantedBy=timers.target