# Fichiers/~ Polytech/NVIDIA/TD-9-11-22/Arthur/dump-database.timer
[Unit]
Description=10min

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target