# Fichiers/~ Polytech/NVIDIA/TD-9-11-22/Arthur/dump-database.service
[Unit]
Description=Sauvegarde_base_donnée 

[Service]
ExecStart=/home/arthur/dump-database.sh

[Install]
WantedBy=multi-user.target