# Vitis/NVIDIA/TD-9-11-22/Tanguy/dump-database.service
[Service]
Description= Sauvegarde database
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/TdBDD/dump-database.sh

[Install]
WantedBy=multi-user.target